# scraper-scheduler.service
[Unit]
Description=Daily Scraper Scheduler Service
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/the-audition-connection
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node /home/ubuntu/the-audition-connection/scheduler.js
StandardOutput=append:/home/ubuntu/the-audition-connection/logs/service.log
StandardError=append:/home/ubuntu/the-audition-connection/logs/service-error.log

[Install]
WantedBy=multi-user.target

